# with-init-system/supervisord.conf
[supervisord]
nodaemon=true

[program:sleep1]
autorestart=false
command=sleep 120

[program:sleep2]
autorestart=false
command=sleep 120
